# service-discovery/service/supervisor.d/app.ini
[program:app]
directory=/app
command=java
    -javaagent:/app/opentelemetry-javaagent.jar 
    -Dotel.service.name=%(ENV_SERVICE_NAME)s
    -Dotel.traces.sampler=parentbased_traceidratio
    -Dotel.traces.sampler.arg=0.1
    -jar /app/main.jar --spring.application.name=%(ENV_SERVICE_NAME)s --management.otlp.metrics.export.resourceAttributes.service.name=%(ENV_SERVICE_NAME)s --management.otlp.metrics.export.resourceAttributes.service.instance.id=%(host_node_name)s --management.metrics.tags.env=%(ENV_SERVICE_ENV)s
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0