# service-discovery/service/supervisor.d/consul.ini
[program:consul]
directory=/app
command=/app/consul agent -data-dir /app/consul-data -node %(ENV_SERVICE_NAME)s-%(ENV_HOSTNAME)s %(ENV_CONSUL_ARGS)s
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0