# spring-integration/app/supervisor.d/app-8081.ini
[program:app-8081]
directory=/app
command=java -jar /app/app.jar
environment= server.port="8081"
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0