# service-discovery/client/supervisor.d/otelcontribcol.ini
[program:otelcontribcol]
directory=/app
command=./otelcol-contrib --config /etc/otelcol-contrib/extensions.yaml --config /etc/otelcol-contrib/exporters.yaml --config /etc/otelcol-contrib/pipeline.agent.yaml --config /etc/otelcol-contrib/pipeline.app.yaml
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0